# units/nd-shutdown.service
[Unit]
Description=Shutdown through network drive.

[Service]
ExecStart=/usr/local/bin/nd-shutdown
