# units/nd-sync.service
[Unit]
Description=Sync network drive to github repository.

[Service]
ExecStart=/usr/local/bin/nd-sync
